# crontab
# Run email_on_update.sh every hour
00 * * * * zsh ~/git/checkUSCIS/email_on_update.sh